# tornado_proj.ini
[group:tornado_proj]
programs=tornado_proj_1,tornado_proj_2

[program:tornado_proj_1]
command={virtual_env}/bin/python main.py --port=8001
directory={prog_dir}/tornado_proj/
user={user}
autorestart=true
redirect_stderr=true
stdout_logfile={prog_dir}/tornado_proj/logs/tornado_proj_1.log
loglevel=debug

[program:tornado_proj_2]
command={virtual_env}/bin/python main.py --port=8002
directory={prog_dir}/tornado_proj/
user={user}
autorestart=true
redirect_stderr=true
stdout_logfile={prog_dir}/tornado_proj/logs/tornado_proj_2.log
loglevel=debug